# deploy/rtk-uart4-inject.service
[Unit]
Description=RTK UART4 Direct Injection Service
After=network-online.target
Wants=network-online.target

[Service]
User=taki
WorkingDirectory=/home/taki/GCS-UmemotoLab
ExecStartPre=/bin/sleep 10
ExecStart=/home/taki/GCS-UmemotoLab/.venv/bin/python /home/taki/GCS-UmemotoLab/rtk_tools/rtk_forwarder_service.py --config /home/taki/GCS-UmemotoLab/config/rtk_forwarder.yml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
